# self-update/self-update@.service
[Unit]
Description=Update %i binary

[Service]
ExecStart=/usr/local/bin/self-update.sh %i

[Install]
WantedBy=multi-user.target