# self-update/self-update@.service
[Unit]
Description=Update the %I binary

[Service]
Type=oneshot
ExecStart=/usr/local/bin/self-update.sh %I